# data/otc-custom-deploy/etc/systemd/system/VA-updates.path
[Unit]
Description=Path watcher for unattended-upgrades

[Path]
PathChanged=/var/log/unattended-upgrades/unattended-upgrades.log

[Install]
WantedBy=multi-user.target
