# data/otc-custom-deploy/appliance-wizard/backend/wizard-server.path
[Unit]
Description=Start wizard server if run flag is created

[Path]
PathExists=/var/appliance-wizard/RUN.FLAG

[Install]
WantedBy=multi-user.target
