# data/otc-custom-deploy/etc/cron.d/openthinclient_ldap_backup
# perform ldap backups for the openthinclient software suite
0 2     *  * *  root /usr/local/sbin/openthinclient-ldapbackup > /dev/null 2>&1
