# data/otc-custom-deploy/etc/xdg/autostart/caja-Desktop-resize-Fix.desktop
[Desktop Entry]
Type=Application
Name=caja-desktop-resize-fix
Exec=/usr/local/bin/caja-Desktop-resize-Fix