# data/otc-custom-deploy/etc/xdg/autostart/openthinclient-caja-desktop-fix.desktop
[Desktop Entry]
Type=Application
Name=openthinclient-caja-desktop-fix
Exec=/usr/local/bin/openthinclient-caja-desktop-fix